# contrib/init/vadercoind.conf
description "Vadercoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env VADERCOIND_BIN="/usr/bin/vadercoind"
env VADERCOIND_USER="vadercoin"
env VADERCOIND_GROUP="vadercoin"
env VADERCOIND_PIDDIR="/var/run/vadercoind"
# upstart can't handle variables constructed with other variables
env VADERCOIND_PIDFILE="/var/run/vadercoind/vadercoind.pid"
env VADERCOIND_CONFIGFILE="/etc/vadercoin/vadercoin.conf"
env VADERCOIND_DATADIR="/var/lib/vadercoind"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # vadercoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$VADERCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run vadercoind."
        echo "The setting must appear in $VADERCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Vadercoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$VADERCOIND_PIDDIR"
    chmod 0755 "$VADERCOIND_PIDDIR"
    chown $VADERCOIND_USER:$VADERCOIND_GROUP "$VADERCOIND_PIDDIR"
    chown $VADERCOIND_USER:$VADERCOIND_GROUP "$VADERCOIND_CONFIGFILE"
    chmod 0660 "$VADERCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$VADERCOIND_PIDFILE" \
    --chuid $VADERCOIND_USER:$VADERCOIND_GROUP \
    --exec "$VADERCOIND_BIN" \
    -- \
    -pid="$VADERCOIND_PIDFILE" \
    -conf="$VADERCOIND_CONFIGFILE" \
    -datadir="$VADERCOIND_DATADIR" \
    -disablewallet \
    -daemon

